# packaging/dusk-checker/systemd/serviceradar-dusk-checker.service
[Unit]
Description=ServiceRadar Dusk Node Checker
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/dusk-checker -config /etc/serviceradar/checkers/dusk.json
Restart=always
RestartSec=10
LimitNPROC=512
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
