# serviceradar-snmp-checker_1.0.19/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable service
systemctl stop serviceradar-snmp-checker
systemctl disable serviceradar-snmp-checker

exit 0
