# packaging/core/systemd/serviceradar-core.service
[Unit]
Description=ServiceRadar Core Service
After=network.target

[Service]
Type=simple
User=serviceradar
ExecStart=/usr/local/bin/serviceradar-core -config /etc/serviceradar/core.json
Restart=always
RestartSec=10
TimeoutStopSec=20
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
